# root/etc/systemd/system/ntpclient.timer
[Unit]
Description=Run foo weekly and on boot

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target